# init/vallumd.openrc
#!/sbin/openrc-run

description="vallumd"
command="/usr/sbin/vallumd"

command_background="yes"
pidfile="/run/${SVCNAME}.pid"
command_args="$VALLUMD_OPTIONS"
