# init/vallumd.upstart
# vallumd

author "Stijn Tintel <stijn@linux-ipv6.be>"

start on runlevel [2345]
stop on runlevel [!2345]

script
    test -f /etc/default/vallumd && . /etc/default/vallumd
    test -f /etc/sysconfig/vallumd && . /etc/sysconfig/vallumd

    exec /usr/sbin/vallumd $VALLUMD_OPTIONS
end script
